# rfcomm.service
[Unit]
Description=Bluetooth serial service
After=ser2bt-e-paper-status.service
Requires=bluetooth.service

[Service]
Type=simple
KillSignal=SIGINT
ExecStart=/usr/bin/rfcomm watch hci0 1 setsid /sbin/agetty -L rfcomm0 9600 xterm -a pi

[Install]
WantedBy=bluetooth.service
