# rfcomm.service
[Unit]
Description=Bluetooth Serial Service
#Before=ser2bt-e-paper-status.service
Requires=bluetooth.target

[Service]
Type=simple
KillSignal=SIGINT
ExecStart=/usr/bin/rfcomm watch hci0 1 setsid /sbin/agetty -L rfcomm0 9600 xterm -a pi

[Install]
RequiredBy=bluetooth.target
